# tools/quick-scan-mgmt/example/systemd-timer.timer
[Unit]
Description=Run mdatp quick scan daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
